# cron/cleanbackup.crontab
# clean up any backups not used in the past 7 days out of the backup directory
# set backupdir? to know what is the backup folder
# https://serverfault.com/questions/948011/cron-job-to-delete-all-files-older-than-7-days-with-extension-in-a-specific-fold
30 21 * * * find ~/.local/share/nvim/backup -type f -name '*:*' -mtime +6 -delete
